# ansible/roles/xdc-node/templates/xdc-node.service.j2
[Unit]
Description=XDC Network Node
Documentation=https://docs.xdc.community
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ xdc_user | default('root') }}
Group={{ xdc_group | default('root') }}
WorkingDirectory={{ xdc_data_dir | default('/root/XDC-Node') }}

Environment="XDC_DATA_DIR={{ xdc_data_dir | default('/root/XDC-Node') }}"
Environment="XDC_LOG_LEVEL={{ xdc_log_level | default('3') }}"

{% if client == 'XDPoSChain' %}
ExecStart=/usr/local/bin/XDC \
  --datadir {{ xdc_data_dir | default('/root/XDC-Node') }} \
  --port {{ p2p_port | default(30303) }} \
  --rpc \
  --rpcaddr 0.0.0.0 \
  --rpcport {{ rpc_port | default(8545) }} \
  --rpcapi {{ rpc_api | default('eth,net,web3,txpool') }} \
  --rpccorsdomain "*" \
  --rpcvhosts "*" \
  --ws \
  --wsaddr 0.0.0.0 \
  --wsport {{ ws_port | default(8546) }} \
  --wsapi {{ ws_api | default('eth,net,web3') }} \
  --wsorigins "*" \
  --metrics \
  --metrics.addr 0.0.0.0 \
  --metrics.port {{ metrics_port | default(6060) }} \
  --gcmode {{ pruning_mode | default('full') }} \
  --syncmode {{ sync_mode | default('snap') }} \
  --maxpeers {{ max_peers | default(25) }} \
  --cache {{ cache_size | default(4096) }}
{% else %}
ExecStart=/usr/local/bin/erigon-xdc \
  --datadir {{ xdc_data_dir | default('/root/XDC-Node') }} \
  --port {{ p2p_port | default(30303) }} \
  --http \
  --http.addr 0.0.0.0 \
  --http.port {{ rpc_port | default(8545) }} \
  --http.api {{ rpc_api | default('eth,net,web3,txpool') }} \
  --http.corsdomain "*" \
  --http.vhosts "*" \
  --ws \
  --metrics \
  --metrics.addr 0.0.0.0 \
  --metrics.port {{ metrics_port | default(6060) }}
{% endif %}

Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
LimitNOFILE={{ max_open_files | default(65535) }}
LimitNPROC={{ max_processes | default(65535) }}

# Security
NoNewPrivileges=false
ProtectSystem=no
ProtectHome=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xdc-node

[Install]
WantedBy=multi-user.target
